# templates/supervisor.conf.tmpl
[program:{{.ServiceName}}]
command={{.Command}} {{.Args}}
directory={{.WorkingDir}}
autostart=true
autorestart=true
stderr_logfile={{.RootFolder}}/logs/{{.ServiceName}}.err.log
stdout_logfile={{.RootFolder}}/logs/{{.ServiceName}}.out.log
environment={{.Environment}}